# docker_compose/flask/app.ini
[program:app]
command=python /opt/app/app.py
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
