# mongo/mongo.ini
[program:mongo]
command=mongod --config /etc/mongod.conf --noauth
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
